# pcapcopier-v1.7/crontab
* * * * * date >> /home/silentdefense/cron.log 2>&1
